# com.onlime.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.onlime.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/cdiseetheeye/Desktop/Onlime/.venv/bin/python3</string>
        <string>-m</string>
        <string>onlime</string>
        <string>run</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/cdiseetheeye/Desktop/Onlime</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/onlime.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/onlime.err</string>
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
