# past/com.cdiseetheeye.obsidian-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cdiseetheeye.obsidian-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/cdiseetheeye/Desktop/자아/obsidian-sync/main.py</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>StandardOutPath</key>
    <string>/Users/cdiseetheeye/.config/obsidian-sync/sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/cdiseetheeye/.config/obsidian-sync/sync-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/cdiseetheeye/Desktop/자아/obsidian-sync</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/cdiseetheeye/Desktop/자아/obsidian-sync</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
