# 과거(~2026.04.04.)/com.onlime.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.onlime.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/aiparty/Desktop/Onlime/dist/index.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/aiparty/Desktop/Onlime</string>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/aiparty/.onlime/logs/onlime.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/aiparty/.onlime/logs/onlime-error.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
